# modules/500-start-memcached/memcached.conf
[program:memcached]
command = memcached -p 11211 -u nobody -U 0
;command = memcached -u nobody -m 64 -p 11211 -l 127.0.0.1 -U 0
