# modules/500-start-memcached/memcached.ini
[program:memcached]
command = memcached -u memcached -U 0
;command = memcached -u memcached -U 0 -m 64 -c 1024 -P /var/run/memcached/memcached.pid
