# modules/110-conf-memcached/memcached.ini
[program:memcached]
command = memcached -u memcached
;command = memcached -u memcached -m 64 -c 1024 -P /var/run/memcached/memcached.pid
